# src/dashboard/config/dashboard.conf
[program:dashboard.etl.accident_check_reports]
directory = /var/www/datascience-framework/
command = python3 /var/www/datascience-framework/src/dashboard/etl/accident_check_reports.py
numprocs = 1
process_name = %(program_name)s_%(process_num)02d
startretries = 0
autostart = true
autorestart = false
user = ubuntu

[program:dashboard.aggregations.accident_check_vins_checked_lifetime]
directory = /var/www/datascience-framework/
command = python3 /var/www/datascience-framework/src/dashboard/aggregations/accident_check_vins_checked_lifetime.py
numprocs = 1
process_name = %(program_name)s_%(process_num)02d
startretries = 0
autostart = true
autorestart = false
user = ubuntu

[program:dashboard.aggregations.accident_check_vins_clean_lifetime]
directory = /var/www/datascience-framework/
command = python3 /var/www/datascience-framework/src/dashboard/aggregations/accident_check_vins_clean_lifetime.py
numprocs = 1
process_name = %(program_name)s_%(process_num)02d
startretries = 0
autostart = true
autorestart = false
user = ubuntu

[program:dashboard.aggregations.accident_check_vins_reported_lifetime]
directory = /var/www/datascience-framework/
command = python3 /var/www/datascience-framework/src/dashboard/aggregations/accident_check_vins_clean_lifetime.py
numprocs = 1
process_name = %(program_name)s_%(process_num)02d
startretries = 0
autostart = true
autorestart = false
user = ubuntu



